# configs/systemd/system/pipewire.service.d/50-priority.conf
[Service]
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=88
Nice=-15
IOSchedulingClass=realtime
IOSchedulingPriority=2


