# configs/systemd/system/ModemManager.service.d/50-priority.conf
[Service]
# ModemManager RT priority for low-latency telephony
Nice=-10
IOSchedulingClass=1
IOSchedulingPriority=4
LimitRTPRIO=90
LimitRTTIME=infinity
